# packages/ec2mc/ec2mc-0.2.0-py3-none-any.whl/ec2mc/aws_setup_src/user_data/mc_template/crontabs/ec2mc_crontab
SHELL=/bin/bash
PATH=/usr/bin/:/usr/sbin/
@reboot root /home/ec2-user/manage-scripts/startup_script.sh > /dev/null
* * * * * root /home/ec2-user/manage-scripts/connection_listener.sh > /dev/null
