# packages/trimarlib-brio/trimarlib_brio-0.1.0.dev0-py3-none-any.whl/brio/resources/brio.service
[Unit]
Description=Basic Remote IO Application
After=gpio-exporter

[Service]
Type=simple
ExecStart=python3 -m brio

[Install]
WantedBy=default.target
